# roles/django/templates/crontab
*/10 * * * * {{ ansible_user_id }} /usr/bin/redis-cli flushdb
