# roles/django/templates/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
logfile=/tmp/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200
directory={{ proj_name }}
user={{ ansible_user_id }}

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[group:{{ proj_name }}]
programs=gunicorn_{{ proj_name }}

[program:gunicorn_{{ proj_name }}]
command={{ base_path }}/pari_env/bin/gunicorn -c {{ base_path }}/pari/gunicorn_conf.py -p /tmp/gunicorn.pid pari.wsgi
autostart=true
autorestart=true
redirect_stderr=true
